# systemd/realtime-monitor.service
[Unit]
Description=Real-time pCloud File Monitor (SSD optimized)
After=pomme-mount.service network-online.target
Wants=pomme-mount.service network-online.target
RequiresMountsFor=/mnt/ssd /mnt/pomme

[Service]
Type=simple
User=vincent
Group=vincent
WorkingDirectory=/home/vincent
ExecStart=/home/vincent/realtime-monitor.sh
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Ensure the service waits for mounts
ExecStartPre=/bin/sleep 10

[Install]
WantedBy=multi-user.target
